# research/supervisor/test_config.conf
[program:test_python]
command=python -u test.py
directory=/home/warmonger/Dropbox/Study/Diploma/Diploma/research/supervisor
stdout_logfile=/home/warmonger/test_python_output.txt
redirect_stderr=true
