# cx/work/v2hw/etc/ppc860-unified/CM5307PPC/rc.local
# Note:
#     This file is designed to be called from mounthost.sh

echo "Starting rrund"
$HOMEDIR/rrund -d >/dev/null 2>/dev/null &
